# modules/asl3/filesystem/etc/systemd/system/asl3-firstboot.service
[Unit]
Description=AllStarLink First Boot Setup
Before=asterisk.service allmon3.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/asl3-firstboot
RemainAfterExit=yes

[Install]
WantedBy=multi-user.target
